# fork-source/opt/teamviewer/tv_bin/script/teamviewerd.RHEL.conf
description "TeamViewer remote control daemon"

exec /opt/teamviewer/tv_bin/teamviewerd -f

start on (started dbus or started messagebus)

stop  on (stopping dbus or stopping messagebus)

respawn 
respawn limit 5 60

pre-start exec /opt/teamviewer/tv_bin/script/tv-delayed-start.sh
