# iso-build/includes.chroot/etc/xdg/autostart/anydesk_global_tray.desktop
[Desktop Entry]
Name=AnyDesk Tray
GenericName=AnyDeskTray
X-GNOME-FullName=AnyDesk
Exec=/usr/bin/anydesk --tray
Icon=anydesk
Terminal=false
Type=Application
